# run_timer/daily_my_job.service
[Unit]
Description=My job run daily.

[Service]
Type=oneshot
ExecStart=/my/running/script/path

[Install]
WantedBy=daily_my_job.service
